# supervisor-ml.conf
[program:ml]
command=python rest_api.py

directory=/srv/giotto/ml/server/
autostart=true
autorestart=true
environment=PYTHONPATH="$PYTHONPATH:/srv/"
stopasgroup=true
